# crontab
PATH=/usr/local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
0 0 * * * root /app/daily_spiders.sh >> /var/log/cron.log 2>&1
